# crontab/crawl_proxy_api.sh
#!/usr/bin/env bash
cd ~/courtpy
python helper/crawl_proxy_api.py &
