# crontab/crawl_proxy_api.sh
#!/usr/bin/env bash
cd /home/apps/proxy.spider
python helper/crawl_proxy_api.py &
