# crontab/valid_tmp.sh
#!/usr/bin/env bash
cd /home/apps/proxy.spider
python helper/valid_tmp.py &
