# crontab/valid_drop_to_all.sh
#!/usr/bin/env bash
cd ~/courtpy
python helper/valid_drop_to_all.py &
