# crontab/valid_all_to_drop.sh
#!/usr/bin/env bash
cd ~/courtpy
python helper/valid_all_to_drop.py &
