# crontab/valid_drop_to_all.sh
#!/usr/bin/env bash
cd /home/apps/proxy.spider
python helper/valid_drop_to_all.py &
