# crontab/valid_all_to_drop.sh
#!/usr/bin/env bash
cd /home/apps/proxy.spider
python helper/valid_all_to_drop.py &
