# crontab/valid_tmp.sh
#!/usr/bin/env bash
cd ~/courtpy
python helper/valid_tmp.py &
